# docker/dev/php/supervisor/websockets.conf
[program:websockets]
command=php /app/artisan websockets:serve
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
user=root
stderr_logfile=/var/log/supervisor/websockets.err.log
stdout_logfile=/var/log/supervisor/websockets.out.log